# files/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# Put 5G WLAN interrupts on CPU1
# To find out interrupts run 'cat /sys/kernel/debug/ieee80211/phy0/mwlwifi/info'
echo 2 > /proc/irq/48/smp_affinity

exit 0
